# deploy/rpg-recorder.service
[Unit]
Description=RPG Recorder Discord Bot
Documentation=https://github.com/beisel-it/rpg-recorder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/rpg-recorder
ExecStart=%h/rpg-recorder/venv/bin/python -m bot
EnvironmentFile=%h/rpg-recorder/.env

# Output dirs (override in .env if needed)
Environment=SESSIONS_DIR=%h/rpg-recorder/sessions
Environment=LOG_DIR=%h/.local/share/rpg-recorder/logs

# Restart on crash with short back-off
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Redirect logs to journald  →  journalctl --user -u rpg-recorder -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpg-recorder

# Graceful stop: let ongoing ffmpeg finalization finish (up to 120 s)
TimeoutStopSec=120

[Install]
WantedBy=default.target
